# scripts/crontab
# Edit this file to introduce tasks to be run by cron.
#
# Each task to run has to be defined through a single line
# indicating with different fields when the task will be run
# and what command to run for the task
#
# To define the time you can provide concrete values for
# minute (m), hour (h), day of month (dom), month (mon),
# and day of week (dow) or use '*' in these fields (for 'any').
#
# Notice that tasks will be started based on the cron's system
# daemon's notion of time and timezones.
#
# Output of the crontab jobs (including errors) is sent through
# email to the user the crontab file belongs to (unless redirected).
#
# For example, you can run a backup of all your user accounts
# at 5 a.m every week with:
# 0 5 * * 1 tar -zcf /var/backups/home.tgz /home/
#
# For more information see the manual pages of crontab(5) and cron(8)
#
# m h  dom mon dow   command

# sudo crontab -u www-data -e

# Adjustments may be needed on times

SHELL=/bin/bash
BASE=/opt/hamclock-backend
VENV=/opt/hamclock-backend/venv
VIRTUAL_ENV=/opt/hamclock-backend/venv
HOME=/var/www
XDG_CACHE_HOME=/opt/hamclock-backend/tmp
MPLCONFIGDIR=/opt/hamclock-backend/tmp/mpl
PATH=/opt/hamclock-backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
LC_ALL=C
LANG=C

# get the latest stable and beta versions
57 * * * *     $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/update_versions.pl

44 * * * *     $BASE/scripts/cron-wrapper.sh $VENV/bin/python3 $BASE/scripts/ssn_simple.py
*/5 * * * *    $BASE/scripts/cron-wrapper.sh $VENV/bin/python3 $BASE/scripts/swind_simple.py
0-58/2 * * * * $BASE/scripts/cron-wrapper.sh $VENV/bin/python3 $BASE/scripts/bz_simple.py
*/30 * * * *   $BASE/scripts/cron-wrapper.sh $VENV/bin/python3 $BASE/scripts/flux_simple.py
1-57/4 * * * * $BASE/scripts/cron-wrapper.sh $VENV/bin/python3 $BASE/scripts/xray_simple.py

# the json updates "estimated" multiple times between its regular 3-hour updates. Also, the 3-hour update isn't
# an exact time after the 3rd hour. So for the 3 hour updates, we'll do 10 and 40 (AI says it should be ready
# by 40) and then just do this every hour to catch estimated updates.
10,40 * * * *  $BASE/scripts/cron-wrapper.sh $VENV/bin/python3 $BASE/scripts/kindex_simple.py

20,50 * * * *  $BASE/scripts/cron-wrapper.sh $VENV/bin/python3 $BASE/scripts/dst_simple.py
25 8 * * *     $BASE/scripts/cron-wrapper.sh $VENV/bin/python3 $BASE/scripts/gen_dxpeditions_spots.py

0 1 * * *      $BASE/scripts/cron-wrapper.sh $BASE/scripts/gen_solarflux-history.sh
0 1 * * *      $BASE/scripts/cron-wrapper.sh $BASE/scripts/gen_ssn_history.pl
15 3 * * 0     $BASE/scripts/cron-wrapper.sh $BASE/scripts/update_pota_parks_cache.sh

*/30 * * * *   $BASE/scripts/cron-wrapper.sh $VENV/bin/python3 $BASE/scripts/web15rss_fetch.py

*/2 * * * *    $BASE/scripts/cron-wrapper.sh OWM_REQS_PER_RUN=50 $VENV/bin/python3 $BASE/scripts/update_world_wx.py

5 */8 * * *    $BASE/scripts/cron-wrapper.sh $BASE/scripts/gen_contest-calendar.sh

# Runs AMSAT Active Satellite Filter and builds TLEs
0 */2 * * *    $BASE/scripts/cron-wrapper.sh $BASE/scripts/fetch_tle.sh

# Give NOAA time to breathe
2,17,32,47 * * * * $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/gen_aurora.py
*/30 * * * *   $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/gen_noaaswx.sh

# POTA, SOTA, WWFF
0 3 * * 1      $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/update_wwff_cache.pl
0 4 * * 1,4    $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/update_sota_cache.pl
1-59/2 * * * * $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/gen_onta.pl

1-56/5 * * * * $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/gen_drap.sh

20 2 * * *     $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/gen_cty_wt_mod.sh

*/30 * * * *   $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/update_all_sdo.sh

# maps
10,40 * * * *  $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/kc2g_muf_heatmap.sh
20 */3 * * *   $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/update_cloud_maps.sh
*/5 * * * *    $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/update_drap_maps.sh
3,18,33,48 * * * * $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/update_aurora_maps.sh
0 2,8,14,20 * * *   $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/update_wx_mb_maps.sh

*/30 * * * *   $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/probe_dynamic_endpoints.sh
*/5 * * * *    $BASE/scripts/cron-wrapper.sh /opt/hamclock-backend/scripts/generate_status_json.sh
